# extras/ansible-service-broker.service
[Unit]
Description=Ansible Service Broker
After=network.target

[Service]
User=ansibleservicebroker
Group=ansibleservicebroker
ExecStart=/usr/bin/ansible-service-broker
Restart=on-failure

[Install]
WantedBy=multi-user.target
